# crontab
* * * * * cd /home/ubuntu/contracts && git fetch origin master && git reset --hard FETCH_HEAD
*/5 * * * * /home/ubuntu/contracts/scripts/sync_with_document_cloud.sh  # Syncs DocumentCloud project with local database

0 4 * * * /home/ubuntu/contracts/scripts/backup.sh  # Creates backups

0 10 * * * /home/ubuntu/contracts/scripts/main.sh  # Check for new contracts
0 12 * * * /home/ubuntu/contracts/scripts/main.sh  # and update app.
0 14 * * * /home/ubuntu/contracts/scripts/main.sh
0 16 * * * /home/ubuntu/contracts/scripts/main.sh
0 18 * * * /home/ubuntu/contracts/scripts/main.sh
